# dot_config/systemd/user/sshfs@.service
[Unit]
Description=SSHFS mount: %I
Documentation=man:sshfs(1)
After=network-online.target
Wants=network-online.target
# %I is the unescaped instance: user@host:/remote:/local
# Use sshfs.bash for interactive mount management.

[Service]
Type=simple

ExecStart=%h/bin/sshfs-mount.sh start "%I"
ExecStop=%h/bin/sshfs-mount.sh stop "%I"

# Restart on connection loss (sshfs exits with code 1 after unmounting).
# RestartSec=30 avoids hammering a temporarily unreachable host.
Restart=on-failure
RestartSec=30

# Limit to 5 restarts within 10 minutes; after that the unit enters failed
# state so it does not loop forever when the host is permanently gone.
StartLimitIntervalSec=600
StartLimitBurst=5

[Install]
WantedBy=default.target
