# dot_config/systemd/user/yt-mp3.timer
[Unit]
Description=Run yt-mp3.service weekly on Sundays

[Timer]
OnCalendar=Fri *-*-* 13:00:00
Persistent=true

[Install]
WantedBy=timers.target
